# syncthing@.timer
[Unit]
Description="Starts syncthing@.service after 1 min on boot"

[Timer]
OnBootSec=1min

[Install]
WantedBy=multi-user.target
